# volume_service/audio-hardware.service
[Unit]
Description=AeroJudge PCB 3.6 DAC and amplifier lifecycle
Wants=sound.target
After=sound.target
Before=volume.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/aerojudge-audio-hardware start
ExecStop=/usr/local/sbin/aerojudge-audio-hardware stop
RemainAfterExit=yes
TimeoutStartSec=infinity
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
